# debian/slurm-llnl-setup-wckeys.cron.d
50 * * * * /usr/sbin/slurm-wckeys-setup 2>&1 
