# debian/slurm-llnl-setup-wckeys.cron.d
50 * * * * root /usr/sbin/slurm-wckeys-setup >/dev/null 2>&1
