# debian/slurm-llnl-sync-accounts.cron.d
0 3 * * * root /usr/sbin/slurm-sync-accounts
